# scripts/zakura-dashboard.service
# Always-on Zebra sync observability dashboard for the bench box.
#
# Serves http://<box>:8090/ persistently:
#   /            live charts for any zakurad currently running on the box, plus a
#                run selector to replay every recorded checkpoint-sync-bench run
#   /runs        JSON index of recorded runs under the archive dir
#   /data?run=ID replay of a recorded run (bottleneck verdict recomputed)
#
# The checkpoint-sync benchmark (scripts/checkpoint-sync-bench.sh) records each
# run into $BENCH_HOME/dashboard/runs/<run-id>/, which this service reads.
#
# The service runs the dashboard from a STABLE copy under the archive root rather
# than the bench build clone ($BENCH_HOME/src), because the build clone is left on
# a detached, per-commit checkout and may not hold this script between builds.
#
# One-time install on the self-hosted runner (roman-zakura-3 / label zakura-bench),
# which runs as root:
#   sudo mkdir -p /opt/zakura-bench/dashboard/runs
#   sudo cp scripts/zakura-metrics-dashboard.py /opt/zakura-bench/dashboard/
#   sudo cp scripts/zakura-dashboard.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now zakura-dashboard
#   # open port 8090 to whoever needs to "log in" to analyze runs
# Refresh the dashboard code with the same `cp` of zakura-metrics-dashboard.py
# followed by `sudo systemctl restart zakura-dashboard`.

[Unit]
Description=Zebra sync observability dashboard (live + recorded-run archive)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
# Edit both paths if BENCH_HOME differs from /opt/zakura-bench.
ExecStart=/usr/bin/python3 /opt/zakura-bench/dashboard/zakura-metrics-dashboard.py \
    --archive /opt/zakura-bench/dashboard/runs \
    --host 0.0.0.0 --port 8090
Restart=on-failure
RestartSec=5
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
